# ansible/roles/lampi/templates/blynk_client.conf
[program:blynk_client]
command=node /home/pi/connected-devices/Lampi/blynk/blynk_client.js
environment=HOME="/home/pi"
priority=600
user=pi
autostart=true
autorestart=true