# ansible/roles/lampi/templates/fbcp.conf
[program:fbcp]
command=/usr/local/bin/fbcp
priority=100
autostart=true
autorestart=true
