# ansible/roles/lampi/templates/lamp_ui.conf
[program:lamp_ui]
command={{ repo_directory }}/Lampi/main.py
directory={{ repo_directory }}/Lampi/
environment=HOME="/home/pi"
priority=400
user=pi
autostart=true
autorestart=true

