# ansible/roles/lampi/templates/lamp_service.conf
[program:lamp_service]
command={{ repo_directory }}/Lampi/lamp_service.py
directory={{ repo_directory }}/Lampi/
environment=HOME="/home/pi"
priority=300
user=pi
autostart=true
autorestart=true

