# ansible/roles/lampi/templates/pigpiod.conf
[program:pigpiod]
command=/usr/local/bin/pigpiod -l -g
priority=200
autostart=true
autorestart=true
