# ansible/roles/cloud_web/templates/mqtt-daemon.conf
[program:mqtt-daemon]
command={{ django_dir }}/manage.py mqtt-daemon
directory={{ django_dir }}
priority=400
user=ubuntu
autostart=true
autorestart=true

